# systemd/debate-agent-cleanup.timer
[Unit]
Description=Debate Agent Data Retention Cleanup
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=timers.target