# systemd/debate-agent.service
[Unit]
Description=Debate Agent Multi-System
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=%h/projects/debate-agent

# Environment & Path
Environment="PATH=%h/.local/bin:%h/projects/debate-agent/.venv/bin"
EnvironmentFile=%h/.config/debate-agent/.env

# Startup
ExecStartPre=%h/.local/bin/uv python -c "import litellm, chainlit, duckduckgo_search"
ExecStart=%h/.local/bin/uv run chainlit run src/ui/chainlit_app.py --port 7860 --host 127.0.0.1

# Resilience
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=debate-agent

# Security Hardening (Linux Kernel)
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=%h/projects/debate-agent/logs %h/projects/debate-agent/memory

[Install]
WantedBy=default.target